# crontabAdmin/Applications/Crontab/Crontab/run_dir/code.crontab
# 商品码订单状态
15 * * * * E:\wamp\bin\php\php5.3.10\php Crontab/code/index.php