# crontabAdmin/Applications/Crontab/Crontab/run_dir/yanchu.crontab
# 演出票数据抓取
59 9 * * * F:\wamp\bin\php\php5.3.10\php Crontab/yanchu/index.php