# crontabAdmin/Applications/Crontab/Crontab/cron_dir/companysync.crontab
# 公司数据定时抓取
10 0 * * * D:\wamp\bin\php\php5.3.10\php Crontab/companysync/index.php