# crontabAdmin/Applications/Crontab/Crontab/cron_dir/dongwang.crontab
# 动网场馆数据抓取
30 0 * * * C:\wamp\bin\php\php5.3.10\php Crontab/dongwang/index.php