# crontabAdmin/Applications/Crontab/Crontab/cron_dir/weixinpay.crontab
# 微信支付成功后，对卡充值的操作
* * * * * E:\wamp\bin\php\php5.3.10\php Crontab/weixin/index.php