# crontabAdmin/Applications/Crontab/Crontab/cron_dir/weixinpay.crontab
# 微信支付成功后，对卡充值的操作
* * * * * F:\wamp-php5.5\bin\php\php5.5.12\php Crontab/weixin/index.php