# crontabAdmin/komovie.crontab
# 抠电影影院数据抓取
50 0 * * * D:\wamp\bin\php\php5.3.10\php Crontab/komovie/index.php